# auto-deploy.service
[Unit]
Description=Auto-Deploy Service - Automated Git Repository Deployment (venv)
After=network-online.target
Wants=network-online.target

[Service]
# Run the auto-deploy Python from the repository virtualenv so host paths
# (project directories) are available to the service.
Type=simple
Restart=always
RestartSec=5
ExecStart=/root/prod/auto-deploy-service/.venv/bin/python /root/prod/auto-deploy-service/auto-deploy.py /root/prod/auto-deploy-service/config.yaml
ExecStop=/bin/sh -c '/root/prod/auto-deploy-service/.venv/bin/python -c "import sys; sys.exit(0)" || true'
WorkingDirectory=/root/prod/auto-deploy-service
StandardOutput=journal
StandardError=journal
User=root
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
